# config/templates/flask-api.service
[Unit]
Description=Flask API for Docker Stats Monitoring
Documentation=https://docs.dusk.network
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
# Run as current user (not root) - will be replaced by installer
User=__USER__
Group=__USER__
WorkingDirectory=__BASE_DIR__
ExecStart=/usr/bin/python3 __BASE_DIR__/app.py
Environment=FLASK_API_PORT=5000
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=/var/crash

# Resource limits
LimitNOFILE=1024
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flask-api

# Environment variables
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

